# sh/lurk_cleanup.timer
[Unit]
Description=Runs cleanup every 48 hours at 4:00 AM

[Timer]
OnCalendar=*-*-* 04:00:00
Persistent=true

[Install]
WantedBy=timers.target