# sh/lurk-cleanup.service
[Unit]
Description=Delete old video recordings

[Service]
Type=onshot
ExecStart=/home/lurk/LurkAlert/sh/cleanup_recordings.sh